# svoydom-crm-bot.service.example
[Unit]
Description=SvoyDom CRM Bot backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/svoydom-crm-bot
EnvironmentFile=/opt/svoydom-crm-bot/.env
ExecStart=/opt/svoydom-crm-bot/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
